# api.service
[Unit]
Description=API server
Alias=api.service

[Service]
ExecStart=/home/nurlybek/service-api/api

[Install]
WantedBy=multi-user.target
